# configs/hue-recovery.service
[Unit]
Description=HueMedia Recovery Server (port 3000)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hue-media/recovery
ExecStart=/usr/bin/python3 /opt/hue-media/recovery/server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
