# configs/hue-auto-update.service
[Unit]
Description=HueMedia Auto-Update (boot-time)
After=network-online.target
Wants=network-online.target
Before=hue-player.service hue-api.service hue-bluetooth.service

[Service]
Type=oneshot
ExecStart=/opt/hue-media/scripts/auto-update.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
